# phin-bridge-1dot1/files/etc/uci-defaults/99_lock-serial-console
#!/bin/sh

uci -q batch <<-EOT
        set system.@system[0].ttylogin='1'
EOT
