# phin-bridge-1dot1/files/etc/uci-defaults/99_remove-wireless-ap
#!/bin/sh

uci -q batch <<-EOT
    delete wireless.@wifi-iface[0]
EOT
