# playground/etc/init.d/playgroundPsr
#!/bin/sh /etc/rc.common

START=95

boot() {
	if ! [ -f /opt/playground/default.psr ]
    	then
        	node /opt/playground/bin/generateDefaultPSR.js
    	fi
}
